# debian/lxcfs.postinst
#!/bin/sh

set -e

case "$1" in
  configure)
    if test -n "$2"; then # upgrade
        dh_action="reload-or-try-restart";
    else
        dh_action="start"
    fi

    unit=lxcfs.service
    if systemctl -q is-enabled "$unit"; then
        systemctl --system daemon-reload >/dev/null || true
        deb-systemd-invoke $dh_action "$unit"
    fi
  ;;

esac

#DEBHELPER#

exit 0
